# tradingagents-cn/rootfs/etc/services.d/tradingagents/finish
#!/usr/bin/with-contenv bashio

bashio::log.info "TradingAgents-CN service stopped"