# tradingagents-cn/rootfs/etc/services.d/tradingagents/finish
#!/usr/bin/with-contenv bash
# 可选：容器退出时清理动作
exit 0