# tradingagents-cn/rootfs/etc/services.d/tradingagents/finish
#!/usr/bin/with-contenv bash
# 容器停止时执行的清理动作（可选）
echo "TradingAgents-CN container is stopping..."
