# tradingagents-cn/rootfs/etc/services.d/tradingagents/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting TradingAgents-CN..."

# 等待网络就绪
bashio::net.wait_for 8501

# 激活虚拟环境
source /opt/venv/bin/activate

# 切换到工作目录
cd /opt/tradingagents

# 启动应用
exec /usr/bin/start-tradingagents.sh