# tradingagents-cn/rootfs/etc/services.d/tradingagents/run
#!/usr/bin/with-contenv bash
# 启动 TradingAgents-CN 服务
exec /usr/bin/start-tradingagents.sh
