# tradingagents-cn/rootfs/etc/services.d/tradingagents/run
#!/usr/bin/with-contenv bash
set -euo pipefail
cd /app


# 若存在常见入口，优先使用；否则找首个 .py 尝试用 streamlit 启动
if [ -f app.py ]; then
exec streamlit run app.py --server.port 8501 --server.address 0.0.0.0
elif [ -f Home.py ]; then
exec streamlit run Home.py --server.port 8501 --server.address 0.0.0.0
else
first_py=$(find . -maxdepth 2 -name '*.py' | head -n 1)
if [ -n "$first_py" ]; then
exec streamlit run "$first_py" --server.port 8501 --server.address 0.0.0.0
else
echo "No entrypoint found (app.py/Home.py). Please adjust run script." >&2
sleep 5
exit 1
fi
fi