# yesplaymusic/rootfs/etc/s6-overlay/s6-rc.d/yesplaymusic/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: YesPlayMusic
# Simple startup script for Home Assistant integration
# ==============================================================================

declare netease_api_url
declare port

# Get configuration from Home Assistant
netease_api_url=$(bashio::config 'netease_api_url')
port=$(bashio::config 'port')

bashio::log.info "Starting YesPlayMusic..."
bashio::log.info "NetEase API URL: ${netease_api_url}"
bashio::log.info "Port: ${port}"

# Export environment variables
export VUE_APP_NETEASE_API_URL="${netease_api_url}"
export PORT="${port}"
export HOST="0.0.0.0"

# Create data directory if it doesn't exist
mkdir -p /data

# Change to app directory and start
cd /app || exit 1
exec node server.js