# yesplaymusic/rootfs/etc/s6-overlay/s6-rc.d/yesplaymusic/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: YesPlayMusic
# ==============================================================================

bashio::log.info "Starting YesPlayMusic Home Assistant Addon..."

# Get configuration from Home Assistant
NETEASE_API_URL=$(bashio::config 'netease_api_url')
WEB_PORT=$(bashio::config 'port')

# Set environment variables
export VUE_APP_NETEASE_API_URL="${NETEASE_API_URL}"
export PORT="${WEB_PORT}"
export HOST="0.0.0.0"

bashio::log.info "Configuration:"
bashio::log.info "  API URL: ${NETEASE_API_URL}"
bashio::log.info "  Port: ${WEB_PORT}"
bashio::log.info "  Host: ${HOST}"

# Change to app directory
cd /app || bashio::exit.nok "Could not change to /app directory"

# Start the application
bashio::log.info "Starting YesPlayMusic server..."
exec ./start.sh