# yesplaymusic/rootfs/etc/s6-overlay/s6-rc.d/yesplaymusic/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: YesPlayMusic
# Runs YesPlayMusic
# ==============================================================================

declare netease_api_url
declare port

# Get configuration
netease_api_url=$(bashio::config 'netease_api_url')
port=$(bashio::config 'port')

bashio::log.info "Starting YesPlayMusic..."
bashio::log.info "NetEase API URL: ${netease_api_url}"
bashio::log.info "YesPlayMusic will be available on port: ${port}"

# Set environment variables
export VUE_APP_NETEASE_API_URL="${netease_api_url}"
export PORT="${port}"
export HOST="0.0.0.0"

# Create settings directory
mkdir -p /data/settings

# Create basic settings file with API URL
cat > /data/settings/settings.json << EOF
{
  "lang": "zh-CN",
  "appearance": "auto",
  "neteaseApiUrl": "${netease_api_url}"
}
EOF

bashio::log.info "Configuration file created at /data/settings/settings.json"

# Change to app directory
cd /app || exit 1

# Check if dist directory exists
if [ ! -d "dist" ]; then
    bashio::log.error "Build directory not found, running build..."
    npm run build || bashio::log.error "Build failed"
fi

# Start YesPlayMusic web server
bashio::log.info "Starting YesPlayMusic web server on ${HOST}:${port}..."
exec node server.js