# alger-music-player/rootfs/services.d/unm/finish
#!/usr/bin/with-contenv bashio

bashio::log.info "UnblockNeteaseMusic service stopped."