# alger-music-player/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/var/run/supervisord.pid
childlogdir=/var/log
user=root

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

# UnblockNeteaseMusic 服务
[program:unm]
command=node app.js -p 3001 -o -s %(ENV_UNM_SOURCE)s
directory=/opt/unm
user=app
autostart=true
autorestart=true
startsecs=10
startretries=3
stderr_logfile=/var/log/unm_error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
stdout_logfile=/var/log/unm_output.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
environment=PORT=3001,HOST=0.0.0.0,NODE_ENV=production,UNM_SOURCE=%(ENV_UNM_SOURCE)s
priority=100

# Alger Music Player 后端
[program:alger-backend]
command=node server.js
directory=/app
user=app
autostart=true
autorestart=true
startsecs=10
startretries=3
stderr_logfile=/var/log/alger_backend_error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
stdout_logfile=/var/log/alger_backend_output.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
environment=PORT=3000,NODE_ENV=production,MUSIC_API_URL=%(ENV_MUSIC_API_URL)s
priority=200
depends_on=unm

# Nginx Web 服务器
[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
startsecs=5
startretries=3
stderr_logfile=/var/log/nginx_error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
stdout_logfile=/var/log/nginx_output.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
priority=300
depends_on=unm,alger-backend

# 日志清理任务（可选）
[program:logrotate]
command=/usr/sbin/crond -f -d 8
autostart=true
autorestart=true
stderr_logfile=/var/log/crond_error.log
stdout_logfile=/var/log/crond_output.log
priority=50