# alger-music-player/rootfs/services.d/alger-music/finish
#!/usr/bin/with-contenv bashio

bashio::log.info "Alger Music Player stopped."