# alger-music-player/rootfs/services.d/alger-music/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting Alger Music Player..."

# 获取配置选项
MUSIC_API_URL=$(bashio::config 'music_api_url')
LOG_LEVEL=$(bashio::config 'log_level')

# 设置环境变量
export MUSIC_API_URL="${MUSIC_API_URL}"
export LOG_LEVEL="${LOG_LEVEL}"

bashio::log.info "Music API URL: ${MUSIC_API_URL}"
bashio::log.info "Log Level: ${LOG_LEVEL}"

# 等待 UNM 服务启动
bashio::log.info "Waiting for UnblockNeteaseMusic service..."
for i in $(seq 1 30); do
    if curl -s http://localhost:3001 >/dev/null 2>&1; then
        bashio::log.info "UnblockNeteaseMusic service is ready!"
        break
    fi
    sleep 1
done

# 启动 nginx
bashio::log.info "Starting nginx..."
nginx -g 'daemon off;' &

# 启动 Node.js 应用
bashio::log.info "Starting Alger Music Player backend..."
cd /app
exec node server.js