# sgcc_electricity/rootfs/etc/services.d/sgcc/run
#!/usr/bin/with-contenv bashio
exec /run.sh
